# files/etc/systemd/system/systemd-time-wait-sync.service
[Service]
DynamicUser=yes
TimeoutStartSec=30s
