# files/etc/systemd/system/basic.target.wants/zfs.target
[Unit]
Description=ZFS startup target

[Install]
WantedBy=basic.target
